# target/linux/qualcommax/base-files/etc/uci-defaults/15_nss_cleanup.sh
#!/bin/sh

rm -f \
	/etc/config/ecm \
	/etc/config/nss \
	/etc/config/pbuf \
	/etc/config/skb_recycler \
	/etc/config/smp_affinity

exit 0
