# package/emortal/ap-c360-radio-mode/files/ap-c360-radio-mode.init
#!/bin/sh /etc/rc.common

START=05

boot() {
	/usr/sbin/ap-c360-radio-mode apply-boot >/dev/null 2>&1
}

start() {
	boot
}
