# target/linux/qualcommbe/ipq95xx/base-files/etc/uci-defaults/16_be7000_wireless
#!/bin/sh

. /lib/functions.sh

[ "$(board_name)" = "xiaomi,be7000" ] || exit 0

/sbin/be7000-wireless-defaults >/dev/null 2>&1 || true

exit 0
