# target/linux/qualcommbe/ipq95xx/base-files/etc/uci-defaults/17_be7000_nfc
#!/bin/sh

. /lib/functions.sh

[ "$(board_name)" = "xiaomi,be7000" ] || exit 0

uci -q batch <<-EOF
	set nfc.nfc='nfc'
	set nfc.nfc.nfc_enable='1'
	set nfc.nfc.mode='url'
	set nfc.nfc.url='auto'
	commit nfc
EOF

/etc/init.d/nfc enable

exit 0
