# target/linux/qualcommbe/ipq95xx/base-files/etc/uci-defaults/99_be7000_boot_success
#!/bin/sh

. /lib/functions.sh

[ "$(board_name)" = "xiaomi,be7000" ] || exit 0

/etc/init.d/be7000_boot_success enable
/etc/init.d/be7000_boot_success boot

exit 0
